# scr/blockade-hosts-update.service
[Unit]
Description=Update your hosts file without having to worry

[Service]
Type=oneshot
ExecStart=/usr/bin/blockade --update

[Install]
WantedBy=multi-user.target
